# root/etc/services.d/ipfs/run
#!/usr/bin/with-contenv bash

HOME=/data exec \
    s6-setuidgid abc ipfs daemon
